# .script/systemd/fan_control.timer
[Unit]
Description=Fan Control

[Timer]
OnCalendar=daily
Unit=fan_control.service
Persistent=true

[Install]
WantedBy=timers.target
