# .script/systemd/service/duckdns.timer
[Unit]
Description=Esegue DuckDNS ogni ora

[Timer]
OnCalendar=hourly
Unit=ducknds.service
Persistent=true

[Install]
WantedBy=timers.target
