# .script/systemd/service/fan_control.service
[Unit]
Description=Fan Control

[Service]
Type=oneshot
ExecStart=/home/odroid/./.script/fan_control

[Install]
WantedBy=multi-user.target
